# ubuntu-target2/supervisor.d/distcc.conf
[program:distcc]
command=/usr/bin/distccd --no-detach --log-level=info --listen 0.0.0.0 --allow 0.0.0.0/0
autostart=true
autorestart=true
user=distcc