# ubuntu-target2/supervisor.d/distcc.conf
[program:distcc]
command=/usr/bin/distccd --daemon --no-detach --log-stderr --log-level=info --listen 0.0.0.0 --allow 0.0.0.0/0 --enable-tcp-insecure
autostart=true
autorestart=true
user=distcc
