# ubuntu-target1/supervisor.d/cupsd.conf
[program:cupsd]
command=/usr/sbin/cupsd -c /etc/cups/cupsd.conf
user=root
autostart=true
startsecs=0
autorestart=false
