# ubuntu-target1/supervisor.d/proftpd.conf
[program:proftpd]
command=/opt/proftpd/sbin/proftpd
autostart=true
user=root
startsecs=0
autorestart=false
