# ubuntu-target1/supervisor.d/startup.conf
[program:start.sh]
command=/start.sh
user=root
autorestart=false
