# ubuntu-target1/supervisor.d/unrealircd.conf
[program:unrealircd]
command=/opt/unrealircd/Unreal3.2/unreal start
autostart=true
autorestart=true
user=root
