# ubuntu-target1/supervisor.d/xinetd.conf
[program:xinetd]
command=/usr/sbin/xinetd -dontfork
autorestart=false
user=root
