# ubuntu-target1/supervisor.d/xinetd.conf
[program:xinetd]
command=/usr/sbin/xinetd -dontfork
autostart=true
autorestart=true
user=root
